# service/deepseek-ocr-api.service
[Unit]
Description=DeepSeek-OCR API Service
Documentation=https://github.com/deepseek-ai/DeepSeek-OCR
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/DeepSeek-OCR
Environment="PATH=/root/DeepSeek-OCR/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/root/DeepSeek-OCR/service/deepseek-ocr.env

# Start service
ExecStart=/root/DeepSeek-OCR/.venv/bin/uvicorn api.main:app --host ${API_HOST} --port ${API_PORT} --log-level info

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Resource limits
LimitNOFILE=65535
TimeoutStartSec=300
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepseek-ocr-api

[Install]
WantedBy=multi-user.target
